# autoloop/org-autoloop.service.template
[Unit]
Description={ORG_DISPLAY_NAME} autoloop — sequential plan queue processor
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/repos/{ORG_DIR}/{REPO_NAME}

# Subscription to use. Switch to rolliq if cashbucket usage is exhausted, or vice versa.
# The CLAUDE_CONFIG_DIR determines which subscription pays for each autoloop run.
# To switch: edit this line, then: systemctl --user daemon-reload && systemctl --user restart {ORG}-autoloop.timer
Environment=CLAUDE_CONFIG_DIR=/home/kev/.claude-{ORG_SUBSCRIPTION}

Environment=GH_CONFIG_DIR=/home/kev/.config/gh-{ORG}
Environment=PATH=/home/kev/.local/bin:/usr/local/bin:/usr/bin:/bin

# --- PRE-JUNE-15-2026 config (RUNNER=run-batch.sh, shared subscription pool) ---
# On 15 Jun 2026: set RUNNER=./run-once.sh, uncomment BILLING_DAY + MONTHLY_CREDIT_USD
# Fill in BILLING_DAY = the day of month your subscription renews (1-28)
Environment=RUNNER=./run-batch.sh
#Environment=RUNNER=./run-once.sh
#Environment=BILLING_DAY={BILLING_DAY}
#Environment=MONTHLY_CREDIT_USD={MONTHLY_CREDIT_USD}

Environment=MODEL=
Environment=MAX_ITERS=8
Environment=PER_CALL_BUDGET_USD=0.75
ExecStart=%h/repos/{ORG_DIR}/{REPO_NAME}/queue-runner.sh

# 75 = "paced out / usage window not ready, try later" — not a failure:
SuccessExitStatus=0 75
Nice=10
